# tdx7984.cron
0 1 * * * root systemctl stop tdx7984
0 3 * * sun root systemctl stop tdx7984 ; renew-stop-info.sh
0 5 * * * root systemctl start tdx7984
