# tdx7984.cron
0 */5 * * * root systemctl restart tdx7984
0 3 * * * root systemctl stop tdx7984 ; renew-stop-info.sh
